# share/profile.d/01-path.sh
#!/bin/sh

export PATH="${XDG_BINARY_HOME}:/local/bin:/bin"
